# idol.conf
[program:idol]
user=root
command=/usr/bin/idol
autostart = true
autorestart = true

